# .docker/supervisor/supervisor-elasticsearch.conf
[program:elasticsearch]
command=/opt/docker/scripts/start-elasticsearch.sh
numprocs=1
autostart=true
autorestart=true