# .docker/supervisor/supervisor-mysql.conf
[program:mysqld]
command=/opt/docker/scripts/start-mysql.sh
numprocs=1
autostart=true
autorestart=true